# private_dot_config/private_aconfmgr/files/private_etc/systemd/system/restore-backlight.service
[Unit]
Description="Restore backlight"

[Service]
Type=simple
ExecStart="/usr/bin/light -I"

[Install]
WantedBy=multi-user.target 
